# debian/postrm
#!/bin/sh
systemctl disable timekpr
systemctl stop timekpr
systemctl reload dbus